# scripts/com.aimemory.reindex.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.aimemory.reindex</string>
    <key>ProgramArguments</key>
    <array>
        <string>/usr/bin/curl</string>
        <string>-s</string>
        <string>--max-time</string>
        <string>60</string>
        <string>-X</string>
        <string>POST</string>
        <string>http://127.0.0.1:8420/index</string>
    </array>
    <key>StartInterval</key>
    <integer>300</integer>
    <key>WorkingDirectory</key>
    <string>/Users/mediocre/projects/ai-memory</string>
    <key>StandardOutPath</key>
    <string>/Users/mediocre/projects/ai-memory/logs/reindex.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/mediocre/projects/ai-memory/logs/reindex.err</string>
</dict>
</plist>
